# main/torrentvpn.service
[Unit]
Description=Test for vpn and torrent
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=alex
ExecStart=/home/alex/github/torrentvpn/main/torrentvpn.sh

[Install]
WantedBy=multi-user.target

